# plugins-core/skill-meta-autoimprover/src/com.harvey.autoimprover.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Harvey Autoimprover LaunchAgent template.
  Before loading, replace __HARVEY_HOME__ with your actual path:
    sed "s|__HARVEY_HOME__|$HOME/HARVEY|g" com.harvey.autoimprover.plist > ~/Library/LaunchAgents/com.harvey.autoimprover.plist
    launchctl load ~/Library/LaunchAgents/com.harvey.autoimprover.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.harvey.autoimprover</string>

    <key>ProgramArguments</key>
    <array>
        <string>python3.11</string>
        <string>__HARVEY_HOME__/plugins-core/meta/autoimprover/run_improvements.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__HARVEY_HOME__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HARVEY_HOME</key>
        <string>__HARVEY_HOME__</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartInterval</key>
    <integer>3600</integer>
    <!-- Run every hour. Change to 86400 for daily. -->

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HARVEY_HOME__/tmp/autoimprover.log</string>

    <key>StandardErrorPath</key>
    <string>__HARVEY_HOME__/tmp/autoimprover.error.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
